# rc.local
#!/bin/sh -e
sudo /opt/vc/bin/tvservice -o
sudo /usr/bin/gpio -g mode 18 in
sudo /usr/bin/gpio -g mode 27 out
sudo /usr/bin/gpio -g write 19 1
sudo /usr/bin/gpio -g mode 19 out
sudo /usr/bin/gpio -g write 27 1
sudo  /home/pi/hello
sleep 3
#sudo /usr/bin/dfu-programmer atmega32u4 erase
#sudo /usr/bin/dfu-programmer atmega32u4 flash /home/pi/zero_avr_firmware_std/main.hex
sudo /usr/bin/dfu-programmer atmega32u4 start


exit 0
